# tools/rebase-jayofelony/overlay/etc/systemd/system/pwnghost-rs.service
[Unit]
Description=PWNGHOST-RS - Rust Pwnagotchi Implementation
Documentation=https://github.com/pwnghost-rs/pwnghost-rs
After=network.target
Wants=network.target

[Service]
Type=notify
NotifyAccess=main
ExecStart=/usr/local/bin/pwnghost-rs --config /etc/pwnghost/config.toml
Restart=on-failure
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=3

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/etc/pwnghost /var/log/pwnghost /var/tmp/pwnghost /var/lib/pwnghost
# RuntimeDirectory (not a bare /run/pwnghost in ReadWritePaths=) is
# required: /run is tmpfs, wiped fresh on every boot, so nothing at
# build time can pre-create /run/pwnghost for it to persist. With
# ProtectSystem=strict, listing a ReadWritePaths= entry that doesn't
# exist makes systemd's mount-namespace setup fail outright before the
# process even starts -- confirmed on real hardware: "Failed to set up
# mount namespacing: .../run/pwnghost: No such file or directory",
# "Failed at step NAMESPACE", status=226/NAMESPACE, looping through the
# restart limit and never actually starting. RuntimeDirectory= is
# systemd's own mechanism for exactly this: it creates (and removes)
# /run/pwnghost around each start/stop of this unit.
RuntimeDirectory=pwnghost
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW CAP_SYS_ADMIN CAP_DAC_OVERRIDE CAP_SYS_RESOURCE
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096
MemoryMax=50M
CPUQuota=80%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pwnghost-rs

[Install]
WantedBy=multi-user.target
